# metanode/deploy/metanode-go-sub@.service
[Unit]
Description=Metanode Go Sub (Node %i)
After=metanode-go-master@%i.service
Requires=metanode-go-master@%i.service

[Service]
Type=simple
User=abc
WorkingDirectory=/home/abc/chain-n/mtn-simple-2025/cmd/simple_chain
ExecStart=/home/abc/chain-n/mtn-simple-2025/cmd/simple_chain/simple_chain -config=config-sub-node%i.json

# Environment
Environment=GOTOOLCHAIN=go1.23.5
Environment=GOMEMLIMIT=4GiB
Environment=XAPIAN_BASE_PATH=sample/node%i/data/data/xapian_node

# File descriptors
LimitNOFILE=100000

# Auto-restart
Restart=always
RestartSec=3
StartLimitIntervalSec=300
StartLimitBurst=10

# Logging
StandardOutput=append:/home/abc/chain-n/mtn-consensus/metanode/logs/node_%i/go-sub-stdout.log
StandardError=append:/home/abc/chain-n/mtn-consensus/metanode/logs/node_%i/go-sub-stdout.log

[Install]
WantedBy=multi-user.target
