# metanode/deploy/metanode-rust@.service
[Unit]
Description=Metanode Rust Consensus (Node %i)
After=metanode-go-master@%i.service
After=metanode-go-sub@%i.service

[Service]
Type=simple
User=abc
WorkingDirectory=/home/abc/chain-n/mtn-consensus/metanode
ExecStart=/home/abc/chain-n/mtn-consensus/metanode/target/release/metanode --config-path config/node_%i.toml

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# File descriptors
LimitNOFILE=100000

# Auto-restart
Restart=always
RestartSec=3
StartLimitIntervalSec=300
StartLimitBurst=10

# Logging
StandardOutput=append:/home/abc/chain-n/mtn-consensus/metanode/logs/node_%i/rust.log
StandardError=append:/home/abc/chain-n/mtn-consensus/metanode/logs/node_%i/rust.log

[Install]
WantedBy=multi-user.target
